# systemd_examples/ahmia-crawl-daily.service
[Unit]
Description=Ahmia daily onion crawl
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Use your venv without 'source'
Environment=PATH=/home/juha/raid/github/ahmia-crawler/venv3/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/juha/raid/github/ahmia-crawler
ExecStart=/home/juha/raid/github/ahmia-crawler/run_daily.sh

# One day: 1*24*60*60 = 86400 seconds, -1400 seconds
RuntimeMaxSec=85000

# Safety net for your memory blow-ups
MemoryMax=10G

# Optional, but sensible
KillSignal=SIGTERM
TimeoutStopSec=120
# Log to journald; view via: journalctl --user -u ahmia-crawl-daily -f
# If you prefer a file, uncomment these two lines:
# StandardOutput=append:/home/juha/raid/github/ahmia-crawler/daily.log
# StandardError=append:/home/juha/raid/github/ahmia-crawler/daily.log
